# private_dot_config/private_aconfmgr/files/etc/systemd/private_system/private_systemd-boot-update.service.d/update-efi-variables.conf
[Service]
ExecStart=
ExecStart=bootctl --graceful update
